# v2.0/images/elasticsearch/etc/supervisor/conf.d/elasticsearch.conf
[program:elasticsearch1]
command=gosu elasticsearch /usr/share/elasticsearch/bin/elasticsearch -Des.default.path.home=/usr/share/elasticsearch -Des.default.path.logs=/var/log/elasticsearch -Des.default.path.work=/tmp/elasticsearch -Des.default.path.conf=/etc/elasticsearch -Des.node.name=lambda-architecture