# v2.0/images/elasticsearch/etc/supervisor/conf.d/elasticsearch.conf
[program:elasticsearch1]
command=gosu elasticsearch /usr/share/elasticsearch/bin/elasticsearch