# images/forwarder/etc/supervisor/conf.d/logstash-forwarder.conf
[program:logstash-forwarder]
command=/opt/logstash-forwarder/logstash-forwarder -config /etc/logstash-forwarder/forwarder.json 
stdout_logfile=/tmp/logs/logstash-forwarder.log
redirect_stderr=true